# tools/tracker-reindex.service.in
[Unit]
Description=Tracker database cleanup and reindexing

[Service]
ExecStart=@DEPLOYMENT_PATH@/tracker_reindex.sh
